# delete_file.list
/usr/bin/droid/droid-hal-startup.sh
/usr/lib/systemd/system/droid-hal-init.service
