# sparse/etc/systemd/system/user@.service.d/50-vendor-hwcomposer.conf
[Service]
# stop first in case something else managed to start it
ExecStartPre=/bin/echo Restarting hwcomposer
ExecStartPre=-+/usr/bin/killall vendor.qti.hardware.display.composer-service
ExecStartPre=-+/system/bin/start vendor.qti.hardware.display.composer
ExecStartPre=/bin/echo Done
