# sparse/etc/systemd/system/systemd-ask-password-gui-stop.service.d/50-vendor-hwcomposer.conf
[Service]
ExecStartPre=-/usr/bin/killall systemd-ask-password-gui
# stop first in case something else managed to start it
ExecStartPre=-/usr/bin/killall vendor.qti.hardware.display.composer-service
ExecStartPre=-/system/bin/start vendor.qti.hardware.display.composer
ExecStart=/bin/echo After vendor.hwcomposer restart
